# deploy/eft.service
[Unit]
Description=Extreme Football Throwdown GMod Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=gmod
WorkingDirectory=/home/gmod/server
ExecStart=/home/gmod/server/garrysmod/gamemodes/extremefootballthrowdown/deploy/start.sh
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
